# rootfs/etc/services.d/sockexec/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
s6-setuidgid multistreamer
/usr/bin/sockexec -t0 /tmp/exec.sock

