# rootfs/etc/services.d/dnsmasq/run
#!/bin/execlineb -P

fdmove -c 2 1

dnsmasq -k --log-facility=-
